# brc20/api/opi-lc.service
[Unit]
Description=OPI-LC Service
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=indexer
Group=indexer
WorkingDirectory=/home/indexer/OPI-LC/brc20/api
Environment=NODE_ENV=production
Environment=PORT=3003
ExecStart=/home/indexer/OPI-LC/brc20/api/start_service.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target 